# market_dashboard/scripts/com.mycodex.invest-backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.mycodex.invest-backend</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/xmly/Library/Mobile Documents/com~apple~CloudDocs/MyCodex/market_dashboard/scripts/start_on_phone.sh</string>
	</array>
	<key>WorkingDirectory</key>
		<string>/Users/xmly/Library/Mobile Documents/com~apple~CloudDocs/MyCodex/market_dashboard</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
		<string>/Users/xmly/Library/Mobile Documents/com~apple~CloudDocs/MyCodex/market_dashboard/output/invest-backend.stdout.log</string>
	<key>StandardErrorPath</key>
		<string>/Users/xmly/Library/Mobile Documents/com~apple~CloudDocs/MyCodex/market_dashboard/output/invest-backend.stderr.log</string>
</dict>
</plist>
